# code/scripts/systemd/tichu-client.service
# systemd unit file for Next.js web client (standalone mode)
#
# Edit the paths below to match your deploy target directory.
# Default layout assumes:
#   /files/.www/tichu/built/client/   — deployed client artifacts (Next.js standalone)
#   /files/.www/tichu/env.client      — client environment file
#
# Install:
#   sudo cp tichu-client.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable tichu-client

[Unit]
Description=Tichu Web Client (Next.js)
After=network.target tichu-server.service
Documentation=https://github.com/sifounak/Tichu

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/files/.www/tichu/built/client/packages/client
ExecStart=/usr/bin/node server.js
EnvironmentFile=/files/.www/tichu/env.client
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
