# code/scripts/systemd/tichu-server.service
# systemd unit file for Fastify game server
#
# Edit the paths below to match your deploy target directory.
# Default layout assumes:
#   /files/.www/tichu/built/server/   — deployed server artifacts
#   /files/.www/tichu/.env.server      — server environment file
#
# Install:
#   sudo cp tichu-server.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable tichu-server

[Unit]
Description=Tichu Game Server (Fastify + WebSocket)
After=network.target
Documentation=https://github.com/sifounak/Tichu

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/files/.www/tichu/built/server
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/files/.www/tichu/.env.server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
